# Scripts/com.marcuslai.claudeusage.signing-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marcuslai.claudeusage.signing-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/marcus/Documents/GitHub/personal/claude-usage-widget/Scripts/refresh-signing.sh</string>
    </array>

    <!-- Run at login… -->
    <key>RunAtLoad</key>
    <true/>

    <!-- …and once a day at 10:00 (the script no-ops unless the profile is
         within a few days of expiring, so this is cheap). -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/marcus/Library/Logs/ClaudeUsage-signing-refresh.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/marcus/Library/Logs/ClaudeUsage-signing-refresh.log</string>
</dict>
</plist>
